# systemd/.config/systemd/user/dotfiles-git.service
[Unit]
Description="Auto commit and push dotfile contents"

[Service]
ExecStart=/usr/local/bin/git-push-all.sh

[Install]
WantedBy=multi-user.target
