# systemd/.config/systemd/user/dotfiles-monitor.path
[Unit]
Description="Monitor dotfiles directory"

[Path]
PathModified=/home/wren/dotfiles/
Unit=dotfiles-git.service

[Install]
WantedBy=multi-user.target
